# files/etc/crontabs/adblock
0 06 * * *    /etc/init.d/adblock reload

