# files/etc/crontabs/root
33 03 * * *    /etc/init.d/adblock reload

