# files/etc/uci-defaults/luci
#!/bin/sh

uci -q batch <<-EOT
	set uhttpd.main.redirect_https='0'
	commit
    EOT
